# supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/root/mongodb/bin/mongod --dbpath /root/db

[program:leanote]
command=/bin/bash -c "/root/init.sh && /root/leanote/bin/run.sh"
